# deployment/store.conf
[program:store]
command=/home/green/jper/bin/gunicorn -c /home/green/jper/src/store/deployment/gconf.py store.app:app
user=green
directory=/home/green/jper/src/store
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
